# python_lcd_ui/systemd/filaboss-machine-vision.service
[Unit]
Description=FilaBoss Machine Vision Sender
After=multi-user.target filaboss-ui.service
Wants=filaboss-ui.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=filaboss
Group=filaboss
Environment=FILABOSS_HEADLESS=1
WorkingDirectory=/home/filaboss/FilaBoss/python_lcd_ui
ExecStart=/home/filaboss/FilaBoss/python_lcd_ui/.venv/bin/python3 /home/filaboss/FilaBoss/python_lcd_ui/machine_vision.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
